# roles/download-stack/files/stash.automount
[Unit]
Description=Automount stash

[Automount]
Where=/stash

[Install]
WantedBy=multi-user.target
